# bin/supervisor/zim.ini
[program:zimapi]
directory=/root/code/zim-server/bin/dev
environment=RELEASE=1
command=/root/code/zim-server/bin/dev/zimapi -config config_86.toml -port 1840
redirect_stderr=true
stdout_logfile=/var/log/zim/api.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
startretries=1
autostart=true
autorestart=unexpected
priority=1
user=root

[program:zimbroker]
directory=/root/code/zim-server/bin/dev
environment=RELEASE=1
command=/root/code/zim-server/bin/dev/zimbroker -config config_86.toml
redirect_stderr=true
stdout_logfile=/var/log/zim/broker.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
startretries=1
autostart=true
autorestart=unexpected
priority=1
user=root
